# gpu-dgpu-guard.service
[Unit]
Description=dGPU-mode test dead-man's switch (auto-revert to Hybrid if screen is black)
After=multi-user.target

[Service]
Type=oneshot
# The guard sleeps up to 5 min waiting for /run/dgpu-keep — never time it out.
TimeoutStartSec=infinity
ExecStart=/usr/local/bin/gpu-dgpu-guard.sh

[Install]
WantedBy=multi-user.target
