# gpu-mux-kwin-fix.service
[Unit]
Description=Configure KWin DRM devices based on ASUS GPU MUX mode
Before=plasmalogin.service sddm.service display-manager.service
After=supergfxd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpu-mux-kwin-fix.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
